# script/crontabs/root
0 4 * * 4  sh /data/adb/modules/AnyHosts/script/functions.sh